# units/sickbeard@.service.d/edge.conf
# Run the latest updates from the master branch as well as enable in-app updates
[Service]
SICKBEARD_EDGE_DKROPTS='-e EDGE=1'
