# units/sickbeard@.service.d/sickrage.conf
# Use the sickrage fork instead of the standard sickbeard
[Unit]
Description=Sickrage

[Service]
Environment='DOCKER_IMAGE=needo/sickrage'
